# deployment/systemd/mact-backend.service
[Unit]
Description=MACT Coordination Backend (Unit 1)
Documentation=https://github.com/yourusername/M-ACT
After=network.target

[Service]
Type=simple
User=mact
Group=mact
WorkingDirectory=/opt/mact
Environment="PATH=/opt/mact/.venv/bin"
Environment="PYTHONPATH=/opt/mact"
EnvironmentFile=/opt/mact/deployment/mact-backend.env
ExecStart=/opt/mact/.venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 4 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    backend.app:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mact-backend

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mact/logs

[Install]
WantedBy=multi-user.target
